# etc/init.d/bootlogs
#!/sbin/openrc-run

description="Saves boot messages to /var/log/dmesg"

depend() {
	need localmount
}

start() {
	ebegin "Saving kernel boot messages"

	if [ -f /var/log/dmesg ]; then
		savelog -q -p -c 5 /var/log/dmesg
	else
		touch /var/log/dmesg
		chown root:adm /var/log/dmesg
		chmod 640 /var/log/dmesg
	fi

	if command -v dmesg >/dev/null 2>&1; then
		dmesg -s 524288 > /var/log/dmesg
	elif [ -c /dev/klog ]; then
		dd if=/dev/klog of=/var/log/dmesg &
		sleep 1
		kill $! >/dev/null 2>&1
	fi

	eend $?
}
