# etc/init.d/configfs
#!/sbin/openrc-run

description="Mounts the configfs kernel virtual filesystem"

depend() {
	need sysfs
	use modules
}

start() {
	local mountpoint="/sys/kernel/config"

	if [ ! -d "${mountpoint}" ]; then
		einfo "Configfs support is not available."
		return 0
	fi

	if mountinfo --quiet "${mountpoint}"; then
		return 0
	fi

	ebegin "Mounting configfs virtual filesystem"
	mount -t configfs none "${mountpoint}"
	eend $?
}
