# etc/init.d/console-setup-openrc
#!/sbin/openrc-run

description="Sets the console font and keymap"

depend() {
	need localmount
	provide console-setup
}

start_pre() {
	if [ ! -x /bin/setupcon ]; then
		eerror "setupcon not found, cannot configure console."
		return 1
	fi
}

start() {
	ebegin "Setting up console font and keymap"
	setupcon
	eend $?
}
