# etc/init.d/dev
#!/sbin/openrc-run

description="Manages the /dev filesystem and its sub-filesystems"

depend() {
	need procfs sysfs
}

start() {
	ebegin "Mounting the /dev filesystem"

	if ! mountinfo -q /dev && grep -q devtmpfs /proc/filesystems; then
		mount -t devtmpfs -o mode=0755,nosuid devtmpfs /dev
		if [ $? -ne 0 ]; then
			eend 1 "Failed to mount /dev"
			return 1
		fi
	fi

	ln -sf /proc/self/fd /dev/fd
	ln -sf fd/0 /dev/stdin
	ln -sf fd/1 /dev/stdout
	ln -sf fd/2 /dev/stderr

	if [ ! -d /dev/pts ]; then
		mkdir -p /dev/pts
		chmod 0755 /dev/pts
	fi
	if ! mountinfo -q /dev/pts; then
		mount -t devpts -o noexec,nosuid,gid=5,mode=620 devpts /dev/pts
	fi

	if ! mountinfo -q /dev/shm; then
		mount -t tmpfs -o nodev,nosuid,noexec tmpfs /dev/shm
	fi
	# Set sticky bit for shared memory.
	chmod 1777 /dev/shm

	eend 0
}
