# etc/init.d/hostname
#!/sbin/openrc-run

description="Sets the system hostname"

depend() {
	before net
}

start() {
	# shellcheck disable=SC2154
	local new_hostname="${hostname}"

	if [ -z "${new_hostname}" ] && [ -s /etc/hostname ]; then
		new_hostname=$(head -n 1 /etc/hostname)
	fi

	if [ -n "${new_hostname}" ]; then
		ebegin "Setting hostname to '${new_hostname}'"
		hostname "${new_hostname}"
		eend $?
	else
		einfo "Hostname is not configured in /etc/conf.d/hostname or /etc/hostname"
	fi
}
