# etc/init.d/hwclock
#!/sbin/openrc-run

description="Manages the hardware clock (RTC)"

depend() {
	need localmount
	before shutdown
	keyword -vserver
}

start() {
	einfo "Skipping RTC sync at boot: assumed handled by udev, NTP, or initramfs"
	return 0
}

stop() {
	if ! yesno "${clock_systohc:-yes}"; then
		einfo "Not saving system clock to hardware clock (see /etc/conf.d/hwclock)."
		return 0
	fi

	ebegin "Saving system clock to the hardware clock"
	/sbin/hwclock --systohc ${clock_opts:---utc}
	eend $?
}
