# etc/init.d/localmount
#!/sbin/openrc-run

description="Mounts local filesystems"

depend() {
	need root fsck
}

start() {
	if [ -r /etc/default/locale ]; then
	# shellcheck source=/dev/null
		. /etc/default/locale
		export LANG
	fi

	ebegin "Mounting local filesystems"

	mount -a -t nonfs,nfs4,smbfs,cifs,ncp,ncpfs,coda,ocfs2,gfs,gfs2,ceph \
		-O no_netdev

	eend $?
}
