# etc/init.d/local
#!/sbin/openrc-run

description="Runs local startup and shutdown scripts"

depend() {
	# shellcheck disable=SC2035
	after *
}

start() {
	ebegin "Running local startup scripts"

	if [ -d /etc/local.d ]; then
		for SCRIPT in /etc/local.d/*.start; do
			if [ -f "${SCRIPT}" ] && [ -x "${SCRIPT}" ]; then
				"${SCRIPT}"
			fi
		done
	fi

	if [ -x /etc/rc.local ]; then
		veinfo "Running compatibility script /etc/rc.local"
		/etc/rc.local
	fi

	eend 0
}

stop() {
	ebegin "Running local shutdown scripts"

	if [ -d /etc/local.d ]; then
		find /etc/local.d/ -maxdepth 1 -type f -name "*.stop" | sort -r | \
		while IFS= read -r SCRIPT; do
			if [ -x "${SCRIPT}" ]; then
				"${SCRIPT}"
			fi
		done
	fi

	if [ -x /etc/rc.shutdown ]; then
		veinfo "Running compatibility script /etc/rc.shutdown"
		/etc/rc.shutdown
	fi

	eend 0
}
