# etc/init.d/modules
#!/sbin/openrc-run

description="Loads specified kernel modules"

depend() {
	need dev
	before local
}

start() {
	if [ ! -f /proc/modules ] || [ ! -x /sbin/modprobe ]; then
		return 0
	fi

	ebegin "Loading kernel modules"

	if [ -n "${modules}" ]; then
		for module in ${modules}; do
			veinfo "Loading module ${module} (from conf.d)"
			modprobe "${module}"
		done
	fi

	if [ -f /etc/modules ]; then
		sed -e 's/#.*//' /etc/modules | while read -r module args; do
			if [ -n "${module}" ]; then
				veinfo "Loading module ${module} (from /etc/modules)"
				modprobe "${module}" ${args}
			fi
		done
	fi

	if [ -d /etc/modules-load.d ]; then
		for file in /etc/modules-load.d/*.conf; do
			[ -f "${file}" ] || continue
			sed -e 's/#.*//' "${file}" | while read -r module args; do
				if [ -n "${module}" ]; then
					veinfo "Loading module ${module} (from ${file##*/})"
					modprobe "${module}" ${args}
				fi
			done
		done
	fi

	eend 0 "Failed to load one or more modules"
}
