# etc/init.d/netmount
#!/sbin/openrc-run

description="Mounts network filesystems"

depend() {
	need net
	after localmount
	use nfs
}

start() {
	ebegin "Mounting network filesystems"

	mount -a -t nfs,nfs4,smbfs,cifs,ncp,ncpfs,ceph -O _netdev

	eend $?
}
