# etc/init.d/procfs
#!/sbin/openrc-run

description="Mounts the /proc virtual filesystem"

provide procfs

start() {
	ebegin "Mounting /proc filesystem"

	if ! mountinfo -q /proc; then
		mount -t proc -o "nodev,noexec,nosuid" proc /proc
	fi

	eend $?
}
