# etc/init.d/swap
#!/sbin/openrc-run

description="Activates and deactivates swap partitions and files"

depend() {
	# Swap files may live on other filesystems
	after localmount
	# The root service must run first
	need root
}

start() {
	ebegin "Activating swap partitions and files"
	swapon -a
	eend $?
}

stop() {
	ebegin "Deactivating swap partitions and files"
	swapoff -a
	eend $?
}
