# etc/init.d/swclock
#!/sbin/openrc-run

description="Sets the system clock from the hardware clock"

depend() {
    need localmount
    provide sysclock
    before net
    keyword -vserver
}

start() {
    if ! yesno "${clock_hctosys:-yes}"; then
        einfo "Not setting system clock from hardware clock."
        return 0
    fi

    ebegin "Setting system clock from the hardware clock"
    /usr/sbin/hwclock --hctosys ${clock_opts:---utc}
    eend $?
}
