# etc/init.d/sysctl
#!/sbin/openrc-run

description="Configures kernel parameters from sysctl.conf"

depend() {
	need procfs sysfs
}

start() {
	ebegin "Applying kernel sysctl settings"
	/usr//sbin/sysctl --system
	eend $?
}
