# etc/init.d/sysfs
#!/sbin/openrc-run

description="Mounts the /sys virtual filesystem"

depend() {
	need procfs
}

start() {
	ebegin "Mounting /sys filesystem and its subsystems"

	if ! mountinfo -q /sys; then
		mount -t sysfs -o "nodev,noexec,nosuid" sysfs /sys
	fi

	if [ -d /sys/kernel/security ] && ! mountinfo -q /sys/kernel/security; then
		if grep -q -w securityfs /proc/filesystems; then
			mount -t securityfs securityfs /sys/kernel/security
		fi
	fi

	if [ -d /sys/fs/pstore ] && ! mountinfo -q /sys/fs/pstore; then
		mount -t pstore pstore /sys/fs/pstore
	fi

	local efivarsmnt="/sys/firmware/efi/efivars"
	if [ -d "${efivarsmnt}" ] && ! mountinfo -q "${efivarsmnt}"; then
		mount -t efivarfs efivarfs "${efivarsmnt}"
	fi

	eend $?
}
