# etc/init.d/umountfs
#!/sbin/openrc-run

description="Unmounts all non-essential filesystems"

stop() {
	ebegin "Unmounting filesystems"

	umount -a -r -t "nonfs,nfs4,smbfs,cifs,ncp,ncpfs,coda,ocfs2,gfs,gfs2,ceph"

	eend $?
}
