# etc/init.d/urandom
#!/sbin/openrc-run

description="Saves and restores the kernel's random seed"
random_seed_file="/var/lib/urandom/random-seed"

depend() {
	need localmount
	provide random-seed
}

_save_seed() {
	local pool_bytes
	pool_bytes=$(cat /proc/sys/kernel/random/poolsize 2>/dev/null)
	# Default to 512 bytes if poolsize can't be read
	pool_bytes=${pool_bytes:-4096}

	umask 077
	dd if=/dev/urandom of="${random_seed_file}" bs=$(( (pool_bytes + 7) / 8 )) count=1 status=none
}

start_pre() {
	[ -c /dev/urandom ] || return 0

	checkpath --directory --owner root:root --mode 0700 "$(dirname "${random_seed_file}")"
}

start() {
	ebegin "Initializing random number generator"
	if [ -f "${random_seed_file}" ]; then
		cat "${random_seed_file}" > /dev/urandom
	fi
	if _save_seed; then
		eend 0
	else
		eend 1 "Failed to save new random seed"
	fi
}

stop() {
	[ -c /dev/urandom ] || return 0
	ebegin "Saving random seed"
	if _save_seed; then
		eend 0
	else
		eend 1 "Failed to save random seed"
	fi
}
