# debian/postinst
#!/bin/bash

ln -sv /sbin/openrc-init /sbin/init
ln -av /usr/sbin/poweroff /sbin/poweroff
ln -sv /usr/sbin/reboot /sbin/reboot

rc-update -u
