# debian/postinst
#!/usr/bin/env bash

# -- Exit on errors.


# -- Add link for init PID 1.

ln -sv /sbin/openrc-init /sbin/init

# -- Edit vanilla rc.conf with our default values.

sed -i 's/#rc_parallel="NO"/rc_parallel="YES"/g' /etc/rc.conf
sed -i 's+#rc_shell=/bin/sh+rc_shell=/bin/zsh+g' /etc/rc.conf
sed -i 's/#rc_logger="NO"/rc_logger="YES"/g' /etc/rc.conf
sed -i 's+#rc_log_path="/var/log/rc.log"+rc_log_path="/var/log/rc.log"+g' /etc/rc.conf
sed -i 's+#rc_crashed_start=YES+rc_crashed_start=YES+g' /etc/rc.conf
sed -i 's/#unicode="NO"/unicode="YES"/g' /etc/rc.conf
sed -i 's/#rc_send_sighup="NO"/rc_send_sighup="YES"/g' /etc/rc.conf
sed -i 's/#rc_timeout_stopsec="90"/rc_timeout_stopsec="10"/g' /etc/rc.conf
sed -i 's/#rc_send_sigkill="YES"/rc_send_sigkill="YES"/g' /etc/rc.conf
sed -i 's/#rc_controller_cgroups="YES"/rc_controller_cgroups="YES"/g' /etc/rc.conf
sed -i 's/#rc_cgroup_mode="unified"/rc_cgroup_mode="unified"/g' /etc/rc.conf
sed -i 's/#rc_cgroup_controllers=""/rc_cgroup_controllers="cgroup_root openrc none cpuset cpu cpuacct blkio io memory devices freezer net_cls perf_event net_prio hugetlb pids rdma misc cgroup"/g' /etc/rc.conf

# -- Add ttys.

for i in {2..7}; do
  ln -sv "/etc/init.d/agetty" "/etc/init.d/agetty.tty${i}"
  rc-update add "agetty.tty${i}" default
done

# -- Update services list.

rc-update -u
