# debian/postinst
#!/bin/bash

ln -sv /sbin/openrc-init /sbin/init

rc-update -u
