# debian/preinst
#!/usr/bin/env bash

#############################################################################################################################################################################
#   The license used for this file and its contents is: BSD-3-Clause                                                                                                        #
#                                                                                                                                                                           #
#   Copyright <2025> <Uri Herrera <uri_herrera@nxos.org>>                                                                                                                   #
#                                                                                                                                                                           #
#   Redistribution and use in source and binary forms, with or without modification, are permitted provided that the following conditions are met:                          #
#                                                                                                                                                                           #
#    1. Redistributions of source code must retain the above copyright notice, this list of conditions and the following disclaimer.                                        #
#                                                                                                                                                                           #
#    2. Redistributions in binary form must reproduce the above copyright notice, this list of conditions and the following disclaimer                                      #
#       in the documentation and/or other materials provided with the distribution.                                                                                         #
#                                                                                                                                                                           #
#    3. Neither the name of the copyright holder nor the names of its contributors may be used to endorse or promote products derived from this software                    #
#       without specific prior written permission.                                                                                                                          #
#                                                                                                                                                                           #
#    THIS SOFTWARE IS PROVIDED BY THE COPYRIGHT HOLDERS AND CONTRIBUTORS "AS IS" AND ANY EXPRESS OR IMPLIED WARRANTIES, INCLUDING, BUT NOT LIMITED TO,                      #
#    THE IMPLIED WARRANTIES OF MERCHANTABILITY AND FITNESS FOR A PARTICULAR PURPOSE ARE DISCLAIMED. IN NO EVENT SHALL THE COPYRIGHT HOLDER OR CONTRIBUTORS                  #
#    BE LIABLE FOR ANY DIRECT, INDIRECT, INCIDENTAL, SPECIAL, EXEMPLARY, OR CONSEQUENTIAL DAMAGES (INCLUDING, BUT NOT LIMITED TO, PROCUREMENT OF SUBSTITUTE                 #
#    GOODS OR SERVICES; LOSS OF USE, DATA,   OR PROFITS; OR BUSINESS INTERRUPTION) HOWEVER CAUSED AND ON ANY THEORY OF LIABILITY, WHETHER IN CONTRACT,                      #
#    STRICT LIABILITY, OR TORT (INCLUDING NEGLIGENCE OR OTHERWISE) ARISING IN ANY WAY OUT OF THE USE OF THIS SOFTWARE, EVEN IF ADVISED OF THE POSSIBILITY OF SUCH DAMAGE.   #
#############################################################################################################################################################################



# -- Clear all runlevels to ensure a clean state.

rm -f /etc/runlevels/default/* \
    /etc/runlevels/shutdown/* \
    /etc/runlevels/boot/* \
    /etc/runlevels/sysinit/*


# -- Whitelist of init scripts used in runlevels.

used_services=(
    apparmor-openrc
    acpi-support-openrc
    acpid-openrc
    bcm-drv
    bluetooth-meshd-openrc
    bluetooth-obexd-openrc
    bluetooth-openrc
    bootlogs
    bootmisc
    brightness
    cgroups
    configfs
    cron-openrc
    cryptmount
    dbus-openrc
    debugfs
    dev
    device-mapper
    dmcrypt
    dmeventd
    elogind-openrc
    fsck
    greetd
    hostname
    hwclock
    local
    localmount
    modules
    mountbinds
    netmount
    network-manager
    nohang-desktop-openrc
    ntpsec-openrc
    nvidia-powerd
    plymouth-openrc
    power-profiles-daemon-openrc
    powertop
    procfs
    rapl-power
    root
    rsyslog-openrc
    savecache
    sudo
    swap
    swclock
    sysctl
    sysfs
    thp-settings
    udev
    umountfs
    urandom
    usbmuxd
)


# -- Delete all other files in /etc/init.d.

for f in "${INIT_D_DIR}"/*; do
  fname="$(basename "$f")"
  if [[ ! " ${used_services[*]} " =~ ${fname} ]]; then
    /bin/rm -f "$f"
  fi
done


# -- Create custom runlevel "async" directory if it doesn't exist and stack it to default.

mkdir -p /etc/runlevels/async && rc-update -s add async default


# -- Define the base directory for init scripts.

INIT_D_DIR="/etc/init.d"


# -----------------------------------------------------------------------------
# ## System Initialization (sysinit)
# -----------------------------------------------------------------------------
echo "Populating sysinit runlevel..."

cp "${INIT_D_DIR}"/{apparmor-openrc,bootmisc,cgroups,dev,hostname,modules,procfs,sysctl,sysfs,udev} /etc/runlevels/sysinit/


# -----------------------------------------------------------------------------
# ## Boot-time Services (boot)
# -----------------------------------------------------------------------------
echo "Populating boot runlevel..."

cp "${INIT_D_DIR}"/{bootlogs,brightness,configfs,cryptmount,device-mapper,dmcrypt,dmeventd,fsck,localmount,root,swap,swclock,urandom} /etc/runlevels/boot/


# -----------------------------------------------------------------------------
# ## Main Services (default)
# -----------------------------------------------------------------------------
echo "Populating default runlevel..."

cp "${INIT_D_DIR}"/{acpi-support-openrc,acpid-openrc,dbus-openrc,debugfs,elogind-openrc,greetd,local,mountbinds,netmount,network-manager,nohang-desktop-openrc,ntpsec-openrc,plymouth-openrc,power-profiles-daemon-openrc,powertop,sudo,uuidd} /etc/runlevels/default/


# -----------------------------------------------------------------------------
# ## Asynchronous Services (async)
# -----------------------------------------------------------------------------
echo "Populating async runlevel..."

cp "${INIT_D_DIR}"/{bcm-drv,bluetooth-meshd-openrc,bluetooth-obexd-openrc,bluetooth-openrc,cron-openrc,nvidia-powerd,rapl-power,rsyslog-openrc,thp-settings,usbmuxd} /etc/runlevels/async/


# -----------------------------------------------------------------------------
# ## Shutdown Services (shutdown)
# -----------------------------------------------------------------------------
echo "Populating shutdown runlevel..."

cp "${INIT_D_DIR}"/{brightness,hwclock,savecache,umountfs,urandom} /etc/runlevels/async/

# -----------------------------------------------------------------------------
# ## Rebuild the dependency tree for the boot order to be correct.
# -----------------------------------------------------------------------------
echo "Rebuilding OpenRC dependency cache..."

rc-update -u
