# etc/local.d/cfs_tweaks.start
#!/bin/sh

# tweak cfs scheduler
/usr/bin/set-cfs-zen-tweaks || true
