# etc/local.d/daemonize_schedtool.start
#!/bin/sh

# start schedtool
/usr/bin/daemonize /usr/bin/schedtool || true
